# chef-repo/site-cookbooks/shizai/templates/default/crontab/set_iptables.sh.erb
#!/bin/bash
# iptables の設定用shell

# mysqlをdocker内許可
iptables -I INPUT 1 -p tcp -m tcp --dport 3306 -j ACCEPT -s 172.16.10.0/24
